# keepalivewifi/crontab
*/5 * * * * root /usr/local/bin/keepalivewifi.sh

